# hermes-kiro-gateway.service
[Unit]
Description=Kiro Gateway — Amazon Q Developer API Proxy
After=network.target

[Service]
Type=simple
User=ali
WorkingDirectory=/home/ali/hermes
ExecStart=/home/ali/hermes-agent/venv/bin/python3 /home/ali/hermes/kiro-gateway.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
