# hermes-kilo-forwarder.service
[Unit]
Description=Hermes Kilo Forwarder (api.kilo.ai via Tor)
After=network.target

[Service]
Type=simple
User=ali
WorkingDirectory=/home/ali/hermes
ExecStart=/home/ali/hermes-agent/venv/bin/python3 /home/ali/hermes/kilo_forwarder.py 9001
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
